# com.user.mail_automation.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.user.mail_automation</string>
    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>/Users/4jp/Code/organvm/universal-mail--automation/run_automation.sh</string>
    </array>
    <key>WorkingDirectory</key>
    <string>/Users/4jp/Code/organvm/universal-mail--automation</string>
    <key>StandardOutPath</key>
    <string>/Users/4jp/System/Logs/mail_automation/launchd.stdout.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/4jp/System/Logs/mail_automation/launchd.stderr.log</string>
    <key>StartCalendarInterval</key>
    <dict>
        <key>Hour</key>
        <integer>9</integer>
        <key>Minute</key>
        <integer>0</integer>
    </dict>
    <key>RunAtLoad</key>
    <false/>
</dict>
</plist>
